# crontab
0 10 * * * cd /py_cronjob && /usr/local/bin/python overnight_freeze_notifier.py 2>&1
